# systemd/smartsake-backup.service
[Unit]
Description=SmartSake Database Backup

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/MainCode
ExecStart=/home/pi/MainCode/.venv/bin/python backup_db.py
StandardOutput=journal
StandardError=journal
