# systemd/smartsake-sensors.service
[Unit]
Description=SmartSake Sensor Collector
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartSake
ExecStart=/usr/bin/python3 /home/pi/SmartSake/WriteSensors.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
