# systemd/smartsake-sensors.service
# DEPRECATED — see smartsake.service
# This was the standalone sensor collector from the old split architecture
# (User=pi, /home/pi/MainCode). The canonical service is smartsake.service
# which runs server.py as User=kojitable under /home/kojitable/ClaudeAgents
# and starts the sensor loop in-process.
[Unit]
Description=SmartSake Sensor Collector
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/MainCode
ExecStart=/home/pi/MainCode/.venv/bin/python WriteSensors.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
