# systemd/smartsake-backup.service
[Unit]
Description=SmartSake Database Backup

[Service]
Type=oneshot
User=kojitable
WorkingDirectory=/home/kojitable/ClaudeAgents
ExecStart=/usr/bin/python3 /home/kojitable/ClaudeAgents/backup_db.py
StandardOutput=journal
StandardError=journal
