# systemd/smartsake-server.service
[Unit]
Description=SmartSake Web Server
After=network.target smartsake-sensors.service
Wants=smartsake-sensors.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartSake
ExecStart=/home/pi/SmartSake/.venv/bin/gunicorn -w 2 -b 0.0.0.0:8080 --timeout 60 --access-logfile - server:app
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
