# systemd/smartsake-server.service
# DEPRECATED — see smartsake.service
# This was the gunicorn web server from the old split architecture
# (User=pi, /home/pi/MainCode). The canonical service is smartsake.service
# which runs server.py directly as User=kojitable under /home/kojitable/ClaudeAgents.
[Unit]
Description=SmartSake Web Server
After=network.target smartsake-sensors.service
Wants=smartsake-sensors.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/MainCode
ExecStart=/home/pi/MainCode/.venv/bin/gunicorn -w 2 -b 0.0.0.0:8080 --timeout 60 --access-logfile - server:app
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
