# systemd/smartsake-sensors.service
# FALLBACK — only needed if running server.py under gunicorn (multi-worker).
# The canonical single-process setup (smartsake.service) starts the sensor
# loop as a background thread inside server.py and does not need this.
[Unit]
Description=SmartSake Sensor Collector (standalone)
After=local-fs.target
# Don't run alongside the main service — it already has the sensor loop
Conflicts=smartsake.service

[Service]
Type=simple
User=kojitable
WorkingDirectory=/home/kojitable/ClaudeAgents
ExecStart=/usr/bin/python3 /home/kojitable/ClaudeAgents/WriteSensors.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Resource limits
MemoryMax=128M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
